# etc/service/hosts/run
#!/bin/bash -e

exec /usr/local/bin/docker-host --domain-name=docker /hosts 2>&1

