# etc/service/hosts/run
#!/bin/bash -e

exec /usr/local/bin/docker-host --domain-name=${DOMAIN_NAME:-docker} ${HOSTS_PATH:-/hosts} 2>&1
