# etc/service/hosts/run
#!/bin/bash -e

exec /usr/local/bin/docker-host --domain-name=${DOMAIN_NAME:-docker} /hosts 2>&1

